# bluetooth_esp32/systemd/gant-vocal.service
[Unit]
Description=Inclusiv'Maker - reconnaissance vocale commande gant/pompe
# Attend que le Bluetooth et le son soient prets avant de demarrer, sinon la
# connexion a l'ESP32 ou le peripherique audio peuvent ne pas encore exister.
After=bluetooth.target sound.target network-online.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/louis/inclusive-maker/bluetooth_esp32
# Adresse MAC Bluetooth de l'ESP32 (voir docs/CONTEXTE_PROJET.md).
# GANT_WAKE_WORD peut aussi etre defini ici pour changer le mot declencheur.
Environment=GANT_BT_MAC=30:76:F5:90:97:6A
# Sans ça, Python bufferise sa sortie et les print() n'apparaissent pas en
# temps reel dans "journalctl -f" (pas de TTY sous systemd).
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/louis/inclusive-maker/bluetooth_esp32/venv/bin/python /home/louis/inclusive-maker/bluetooth_esp32/voice_recognition.py --device 0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
